# config/dev.service
[Unit]
Description=Edwin's Central Dev Server

[Service]
Type=notify
ExecStart=/usr/bin/env pnpm --dir %h/.dev start

[Install]
WantedBy=default.target
